# jobs/cfdb_backup/templates/config/backup-db.cron
0 */4 * * * /var/vcap/jobs/cfdb_backup/bin/db_backup
