# jobs/boshdb_backup/templates/config/backup-db.cron
0 */4 * * * /var/vcap/jobs/boshdb_backup/bin/db_backup
